# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

; ── PHP-FPM ───────────────────────────────────────────────────────────────────
[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Nginx ─────────────────────────────────────────────────────────────────────
[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Laravel Horizon (Queue Manager) ──────────────────────────────────────────
; Horizon manages its own worker pool; no need for separate queue:work processes.
; If the entrypoint couldn't reach Redis it drops storage/horizon.disabled and
; switches the queue to database — run a plain queue:work instead of Horizon.
[program:horizon]
command=bash -c "if [ -f /var/www/html/storage/horizon.disabled ]; then exec php /var/www/html/artisan queue:work database --sleep=3 --tries=3 --max-time=3600; else exec php /var/www/html/artisan horizon; fi"
autostart=true
autorestart=true
startretries=5
startsecs=3
stopasgroup=true
killasgroup=true
stopsignal=SIGTERM
stopwaitsecs=120
priority=30
user=www-data
numprocs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Laravel Scheduler ─────────────────────────────────────────────────────────
[program:scheduler]
command=bash -c "while [ true ]; do php /var/www/html/artisan schedule:run --no-interaction & sleep 60; done"
autostart=true
autorestart=true
priority=40
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
